# debian/wiperf.service
[Unit]
Description=Wiperf: UX probe

[Service]
ExecStart=/usr/sbin/wiperf

[Install]
WantedBy=multi-user.target