# debian/wlanpi-wiperf.service
[Unit]
Description=Wiperf: UX probe

[Service]
ExecStart=/opt/wlanpi-wiperf/bin/wiperf-scheduler

[Install]
WantedBy=multi-user.target